# utils/systemd/monerod.service
[Unit]
Description=Monero Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/monerod --config-file /etc/monero/monerod.conf --non-interactive

Restart=always
RestartSec=30

User=monero
Group=monero

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
