# utils/systemd/monerod.service
[Unit]
Description=Monero Full Node
After=network.target

[Service]
User=monero
Group=monero
WorkingDirectory=~

Type=forking
PIDFile=/var/run/monerod.pid

ExecStart=/usr/bin/monerod --config-file /etc/monerod.conf \
    --detach --pidfile /var/run/monerod.pid

[Install]
WantedBy=multi-user.target
